# supervisord.conf
[unix_http_server]

file=/tmp/supervisor.sock

[supervisord]
logfile=/tmp/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/var/run/supervisord.pid
nodaemon=false
minfds=1024
minprocs=200
[rpcinterface:supervisor]

supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock


[program:pedaloServer]
command=/home/pi/PEDALO/pedalo-vm/bin/python3 -m pedalo_server
process_name=%(program_name)s
autostart=true
autorestart=true
